# .config/systemd/user/dunst-resume.service
[Unit]
Description=Resume dunst notifications

[Service]
Type=oneshot
ExecStart=/usr/bin/dunstctl set-paused false

