# .config/systemd/user/dunst-pause.service
[Unit]
Description=Pause dunst notifications

[Service]
Type=oneshot
ExecStart=/usr/bin/dunstctl set-paused true
